# idle.service
[Unit]
Description=Idle service

[Service]
ExecStart=/usr/local/bin/idle.sh
